# configs/releng/airootfs/usr/lib/systemd/system/refresh-package-databases.service
[Unit]
Description=Refresh all package databases

[Service]
Type=oneshot
ExecStart=/bin/bash -c "pacman -Syy"